# settings/systemd/user/xremap.service
[Unit]
Description=xremap key remapper
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.local/share/mise/shims/xremap %h/.config/xremap/config.yml
Restart=on-failure

[Install]
WantedBy=graphical-session.target
